# systemd/backup-config.service
[Unit]
Description=Run Proxmox config backup before shutdown
DefaultDependencies=no
Before=shutdown.target halt.target
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/root/.backup-config.env
ExecStart=/root/backup-config.sh
TimeoutStartSec=600
StandardOutput=append:/var/log/backup-config.log
StandardError=append:/var/log/backup-config.log

[Install]
WantedBy=halt.target
WantedBy=shutdown.target
