# systemd/pa-backup-config.service
[Unit]
Description=Run Proxmox config backup

[Service]
Type=oneshot
EnvironmentFile=/root/.pa-agent.env
ExecStart=/usr/local/bin/pa-backup-config.sh
TimeoutStartSec=600
StandardOutput=append:/var/log/pa-backup-config.log
StandardError=append:/var/log/pa-backup-config.log

[Install]
WantedBy=multi-user.target
